# tools/kiosk/systemd/digital-dash-kiosk.service
[Unit]
Description=Digital Dash Zero-Flash Kiosk
After=plymouth-start.service plymouth-read-write.service network-online.target
Before=graphical.target
Wants=network-online.target
Conflicts=getty@tty1.service display-manager.service

[Service]
Type=simple
User=%i
WorkingDirectory=/digital-dash
Environment=XDG_SESSION_TYPE=wayland
Environment=GTK_THEME=Adwaita:dark
ExecStart=/digital-dash/tools/kiosk/zero-flash-kiosk.sh
Restart=always
RestartSec=3
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=graphical.target