# tools/systemd/digital-dash-bluetooth.service
[Unit]
Description=Digital Dash Bluetooth Service
After=network.target bluetooth.service

[Service]
Type=simple
WorkingDirectory=/digital-dash
ExecStart=/usr/bin/node server/bluetooth-service/server.js
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
