# tools/systemd/digital-dash-vehicle.service
[Unit]
Description=Digital Dash Vehicle Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
ExecStart=__NPM_BIN__ run vehicle-service
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
