# tools/systemd/digital-dash-bluetooth.service
[Unit]
Description=Digital Dash Bluetooth Service
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
ExecStart=__NODE_BIN__ server/bluetooth-service/server.js
Environment=BLUETOOTHCTL_CMD=/usr/bin/bluetoothctl
Environment=PACTL_CMD=/usr/bin/pactl
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
