# tools/systemd/digital-dash-ui.service
[Unit]
Description=Digital Dash UI
After=local-fs.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
Environment=DIGITAL_DASH_UI_HOST=0.0.0.0
Environment=DIGITAL_DASH_UI_PORT=5173
ExecStart=__NPM_BIN__ run ui:serve
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
